# oh/systemd/oh.timer
[Unit]
Description=Regenerate the cache every hour
[Timer]
OnBootSec=15
OnUnitActiveSec=1h
[Install]
WantedBy=timers.target