# arch/systemd/arch-update.timer
[Unit]
Description=Research packages updates every hour
[Timer]
OnBootSec=15
OnUnitActiveSec=1h
[Install]
WantedBy=timers.target