# arch/systemd/arch-update.timer
[Unit]
Description=Check system update for arch
[Timer]
OnBootSec=15
OnUnitActiveSec=1h
[Install]
WantedBy=timers.target
