# arch/systemd/arch.timer
[Unit]
Description=Regenerate the packages cache for arch
[Timer]
OnBootSec=15
OnUnitActiveSec=1h
[Install]
WantedBy=timers.target