# oh.service
[Unit]
Description=Oh service
After=network-online.target

[Service]
Type=oneshot
ExecStart=oh --cache