# oh.timer
[Unit]
Description=Oh Cache Refresh Timer

[Timer]
OnCalendar=hourly
Persistent=true
Unit=oh.service

[Install]
WantedBy=timers.target