# oh.timer
[Unit]
Description=Oh timer

[Timer]
OnStartupSec=15
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target